# deploy/systemd/solana-arb-monitor.service
[Unit]
Description=Solana Arb Bot — monitor stack (Docker Compose)
Documentation=file:///path/to/solana-arb-bot/docs/GO_LIVE.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
Group=docker
WorkingDirectory=/path/to/solana-arb-bot
ExecStartPre=/usr/bin/npm run sync:compose-env
ExecStart=/path/to/solana-arb-bot/scripts/restart-monitor.sh
Restart=always
RestartSec=10
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal
Environment=COMPOSE_SKIP_HEALTH_GATE=0

[Install]
WantedBy=multi-user.target
