# deploy/systemd/solana-arb-monitor.service
[Unit]
Description=Solana Arb Bot — monitor stack (Docker Compose)
Documentation=file:///path/to/solana-arb-bot/docs/GO_LIVE.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=docker
WorkingDirectory=/path/to/solana-arb-bot
ExecStartPre=/usr/bin/npm run sync:compose-env
ExecStart=/path/to/solana-arb-bot/scripts/restart-monitor.sh
Restart=no
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
Environment=COMPOSE_SKIP_HEALTH_GATE=0
Environment=SOPS_AGE_KEY_FILE=/path/to/solana-arb-bot/secrets/.local/sops_age_key

[Install]
WantedBy=multi-user.target
